# systemd/readq3c.service
[Unit]
Description=ReadQ3C
After=network.target

[Service]
ExecStart=/usr/local/bin/readq3c

[Install]
WantedBy=multi-user.target
